# systemd/p1-backup.service
[Unit]
Description=GreenDevCorp P1 Backup Service
Documentation=https://github.com/usuari/GSX-P1
# Storage must be mounted before backup runs
RequiresMountsFor=/mnt/storage

[Service]
Type=oneshot
User=root
ExecStart=/opt/P1/scripts/Week_5_backup/backup.sh
StandardOutput=journal
StandardError=journal
# Keep journal logs for failed runs visible
SyslogIdentifier=p1-backup

[Install]
WantedBy=multi-user.target
